# debian/prerm
#!/bin/sh
# prerm script for geonode
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove|upgrade|deconfigure)
	a2dissite geonode
	a2ensite default
	invoke-rc.d apache2 restart

	rm -rf /var/lib/geonode
	rm -rf /var/www/geonode/wsgi
	rm -rf /var/www/geonode/static
	invoke-rc.d tomcat6 stop
	rm -rf  /var/lib/tomcat6/webapps/geoserver /var/lib/tomcat6/geonetwork
	invoke-rc.d tomcat6 start

	rm -rf /usr/share/geonode/role.sql
	rm -rf /usr/share/geonode/admin.json
	rm -rf /usr/share/geonode/patch*
	rm -rf /usr/bin/geonode

	# The other files that the installer created but we do not want to blow away are:
	#/etc/default/geonode
	#/var/www/geonode/uploaded
	#/var/lib/geoserver/data
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
