# debian/prerm
#!/bin/sh
# prerm script for geonode
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove|upgrade|deconfigure)
	if [ -e /etc/apache2/sites-enabled/0000-geonode ]; then
		invoke-rc.d apache2 stop
		rm -rf /etc/apache2/sites-enabled/0000-geonode /etc/apache2/sites-available/geonode
		if [ ! -e /etc/apache2/sites-enabled/000-default ]; then
			a2ensite default
		fi
		invoke-rc.d apache2 start
	fi

	rm -rf /var/www/geonode 

	# turn off error trapping, one of these may fail
	set +e
	su - postgres -c "dropdb geonode"
	su - postgres -c "dropuser geonode"
	set -e
	# turn it back on

	if [ -e /var/lib/tomcat6/geonetwork ]; then
		invoke-rc.d tomcat6 stop
		rm -rf  /var/lib/tomcat6/webapps/geoserver-geonode-dev /var/lib/tomcat6/geonetwork
		rm -rf /var/lib/tomcat6/webapps/geoserver-geonode-dev.war /var/lib/tomcat6/webapps/geonetwork.war
		invoke-rc.d tomcat6 start
	fi

	rm -rf /usr/share/geonode/role.sql
	rm -rf /usr/bin/geonode
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
